# systemd/bedrock@.service
# Reference template; the installer writes the active copy into /etc/systemd/system
[Unit]
Description=Minecraft Bedrock Server %i
After=network.target

[Service]
WorkingDirectory=/opt/minecraft/%i
User=minecraft
Group=minecraft
Restart=on-failure
RestartSec=10
ExecStart=/usr/bin/screen -DmS mcbe-%i /opt/minecraft/%i/bedrock_server
ExecStop=/usr/bin/screen -S mcbe-%i -p 0 -X stuff "stop\r"

[Install]
WantedBy=multi-user.target
